# init/nagios-helper.conf
# Start nagios-helper
#

description     "Start nagios-helper daemon"

start on startup

task
exec python /usr/local/nagios-helper/nagios-helper